# elements/autoinstall/static/etc/rc.local
#!/bin/bash
#
# rc.local: Place local configuration here.  It is executed by systemd
# during the boot cycle.
#

